# dev-activity.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.activity.watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/stan/dev/dev-activity/.venv/bin/python</string>
    <string>/Users/stan/dev/dev-activity/dev_activity.py</string>
    <string>watch</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/stan/dev/dev-activity</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/stan/dev/dev-activity/watch.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/stan/dev/dev-activity/watch.err</string>
</dict>
</plist>
